# reloadTouchpadDriverOnBoot.service
[Unit]
Description=Reload Touchpad-Driver

[Service]
ExecStart=PATH_TO_YOUR_RELOADTOUCHPAD_SH

[Install]
WantedBy=default.target
